# config/conf/app-uwsgi.ini
[program:uwsgi]

command=uwsgi --ini uwsgi.ini
priority=200
autostart=true
autorestart=true
stopsignal=INT
