# config/conf/app-celery.ini
[program:celery.worker]

directory=../
command=celery -A django_quickstart worker --pool eventlet
					--loglevel info
					--logfile ./log/celery_worker.log

numprocs=1
autostart=true
autorestart=true
stopsignal=INT
priority=100


[program:celery.beat]
directory=../
command=celery -A django_quickstart beat --loglevel info
					--logfile ./log/celery_beat.log

numprocs=1
autostart=true
autorestart=true
stopsignal=INT

priority=120
